# scripts/hub-manager.conf
description "hub-manager"
author "makerlab.mx - info@makerlab.mx"

start on filesystem or runlevel [2345]
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 20 5

# Change home and user to the correct one
script
	iptables -t nat -A PREROUTING -i `ip -o -4 r s to default | grep ^default | awk '{print $5}'` -p tcp --dport 80 -j REDIRECT --to-port 8000
	iptables -t nat -A PREROUTING -i `ip -o -4 r s to default | grep ^default | awk '{print $5}'` -p tcp --dport 443 -j REDIRECT --to-port 8080
	export HOME="/home/odroid"
	exec sudo -u odroid /usr/local/bin/hub-manager --ssl >> /home/odroid/.hub-manager/logs/hub-manager.log 2>&1
end script

post-start script
	echo "INFO: hub-manager started"
end script